# files/crontab


0 2 * * * /usr/local/bin/backup_home.sh >> /var/log/syslog 2>&1